# deploy/wordbot.service
[Unit]
Description=Word of the Day Discord bot
# Wait until the network is up before starting.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/wordbot
ExecStart=/home/admin/wordbot/.venv/bin/python bot.py
# Restart automatically if it ever crashes or the connection drops.
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
